# restic_config/restic_backup.timer
[Unit]
Description=Run restic backup 5 min after start

[Timer]
OnBootSec=5min
Unit=restic_backup.service

[Install]
WantedBy=timers.target